# tools/compose/migworkstation/simpleweb.conf
[program:simpleweb]
command=/mig/simpleweb.py 2222
startretries=9999
autorestart=true
