# tools/compose/migapi/mig-api.conf
[program:mig-api]
command=/go/bin/mig-api
startretries=9999
