# tools/compose/migagenthost2/mig-agent.conf
[program:mig-agent]
command=/go/bin/mig-agent -d
startretries=9999
autorestart=true
