# tools/compose/migscheduler/mig-scheduler.conf
[program:mig-scheduler]
command=/go/bin/mig-scheduler
startretries=9999
autorestart=true
